# infra/systemd/saas-deploy.service
# Unité de déploiement staging (modèle pull). Installée sur la machine de staging :
#   sudo cp infra/systemd/saas-deploy.* /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now saas-deploy.timer
# Suivi : journalctl -u saas-deploy.service
[Unit]
Description=Déploiement staging par pull (images GHCR)
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/srv/saas/app/scripts/deploy-pull.sh
